# files/udev/90_mcloud_device_appium.rules
ACTION=="remove", RUN+="/usr/local/bin/healthcheck"
